# ocr/hermes-ocr.service
[Unit]
Description=Hermes OCR Sync Service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/hermes-ocr/
ExecStart=/opt/hermes-ocr/venv/bin/python /opt/hermes-ocr/ocr_sync.py
Restart=always
RestartSec=10
User=www-data
Group=www-data
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
